# system/Z575/OSes/gentoo/on_baremetal/filesystem_now/gentoo/home/zazdxscf/rsync/rsync_filesfrom.list

/
/usr/src/.config_normal
/usr/src/.config_hardened
/usr/src/.config_normalgit

/usr/src/linux-4.2-rc6/.config
/usr/src/linux-4.4-rc6/.config
/usr/src/linux-4.4-rc7/.config
/usr/src/linux-4.4-rc8/.config
/usr/src/linux-4.1.6-hardened/.config
/usr/src/linux-4.1.6-hardened-r1/.config
/usr/src/linux-4.1.6-hardened-r2/.config
/usr/src/linux-4.1.7-hardened/.config
/usr/src/linux-4.1.7-hardened-r1/.config
/usr/src/linux-4.2.3-hardened/.config
/usr/src/linux-4.2.3-hardened-r1/.config
/usr/src/linux-4.2.3-hardened-r2/.config
/usr/src/linux-4.2.3-hardened-r3/.config
/usr/src/linux-4.2.3-hardened-r4/.config
/usr/src/linux-4.2.3-hardened-r5/.config
/usr/src/linux-4.2.4-hardened/.config
/usr/src/linux-4.2.4-hardened-r2/.config
/usr/src/linux-4.2.5-hardened/.config
/usr/src/linux-4.2.5-hardened-r1/.config
/usr/src/linux-4.2.5-hardened-r2/.config
/usr/src/linux-4.2.6-hardened/.config
/usr/src/linux-4.2.6-hardened-r1/.config
/usr/src/linux-4.2.6-hardened-r3/.config
/usr/src/linux-4.2.6-hardened-r4/.config
/usr/src/linux-4.2.6-hardened-r5/.config
/usr/src/linux-4.2.6-hardened-r6/.config
/usr/src/linux-4.2.6-hardened-r7/.config
/usr/src/linux-4.2.6-hardened-r8/.config
/usr/src/linux-4.2.7-hardened/.config
/usr/src/linux-4.5-rc1/.config
/usr/src/linux-4.5-rc2/.config
/usr/src/linux-4.5-rc3/.config
/usr/src/linux-4.5-rc4/.config
/usr/src/linux-4.5-rc5/.config
/usr/src/linux-4.5-rc6/.config
/usr/src/linux-4.5-rc6/.config
/usr/src/linux-4.5-rc7/.config
/usr/src/linux-4.6-rc1/.config
/usr/src/linux-4.6-rc2/.config
/usr/src/kernel_versions.diff

/usr/src/linux

/etc/init.d/
/etc/default/
/etc/security/
/etc/env.d/
/etc/fstab
/etc/host.conf
/etc/hosts
/etc/inittab
/etc/locale.gen
/etc/rc.conf
/etc/grub.d/40_custom
/etc/kernels/
/etc/metalog.conf
/etc/ntp.conf
/etc/sudoers.d/wheel
/etc/timezone
/etc/updatedb.conf
/etc/vim/
/etc/X11/xorg.conf.d
/etc/xdg/
/etc/portage/
/etc/conf.d/
/etc/ssh/
/etc/portage/repos.conf/gentoo.conf
/etc/resolv.conf
/etc/lvm/





/usr/portage/metadata/layout.conf

/ccache/ccache.conf
/home/zazdxscf/.ccache/ccache.conf


/var/lib/gentoo/news/

/var/lib/portage/

/var/lib/ntp/ntp.drift

/home/zazdxscf/build/1nonpkgs/hosts/
/home/zazdxscf/build/1nonpkgs/rust/
/home/zazdxscf/build/1nonpkgs/rust.vim/
/home/zazdxscf/build/1nonpkgs/cargo/

/home/zazdxscf/build/1nonpkgs/uBlock/go
/home/zazdxscf/build/1nonpkgs/uMatrix/go
/home/zazdxscf/build/1nonpkgs/https-everywhere/go
/home/zazdxscf/build/1nonpkgs/firestorageplus/go
/home/zazdxscf/build/1nonpkgs/footab-code/



/usr/src/compare



/home/zazdxscf/.pki
/home/zazdxscf/.thumbnails
/home/zazdxscf/.cache
/home/zazdxscf/.cargo
/home/zazdxscf/.links



/etc/portage/profile/package.use.mask

.cargo

