# system/Z575/OSes/gentoo/on_baremetal/filesystem_now/gentoo/etc/local.d/000_grsec_lock.start
sysctl -w kernel.grsecurity.grsec_lock=1

