# system/Z575/OSes/gentoo/on_baremetal/filesystem_now/gentoo/etc/local.d/100_delay_writes.start
echo 60000 > /proc/sys/vm/dirty_writeback_centisecs
echo 360000 > /proc/sys/vm/dirty_expire_centisecs

