# system/Z575/OSes/gentoo/on_baremetal/filesystem_now/gentoo/etc/local.d/110_thermal_radeon.start
echo "Before: `cat /sys/class/drm/card0/device/power_dpm_state`"
echo 'thermal' > /sys/class/drm/card0/device/power_dpm_state
echo "After : `cat /sys/class/drm/card0/device/power_dpm_state`"

