# system/Z575/OSes/gentoo/on_baremetal/filesystem_now/gentoo/etc/local.d/120_cpu_governor_set.start
#!/bin/bash

#cpupower frequency-set -g conservative #this will play bad with rust! sticks at 800Mhz while 100% cpu is used by rustc
cpupower frequency-set -g ondemand

