# system/Z575/OSes/gentoo/on_baremetal/filesystem_now/gentoo/etc/local.d/99999-uptime.start
exec 1>/dev/console
exec 2>/dev/console

cat /proc/uptime
uptime
