# system/Z575/OSes/manjaro2/on_baremetal/filesystem_now/manjaro/home/z/rsync/rsync_filesfrom.list
/home/z/.config/gtk-2.0/
/home/z/.config/mc
/root/.config/mc
/home/z/.config/vlc/
/home/z/.config/xfce4/panel/
/home/z/.config/xfce4/xfconf/
/home/z/.config/xfce4/terminal/terminalrc
/home/z/.config/chromium-flags.conf
/root/.vim
/home/z/.vim/
/root/bin
/home/z/bin
/home/z/.bash_profile
/home/z/.bashrc
/home/z/.extend.bashrc
/home/z/.second.bashrc
/home/z/.inputrc
/home/z/.less
/home/z/.vimrc
/home/z/.ccache/ccache.conf
/root/.bash_profile
/root/.bashrc
/root/.extend.bashrc
/root/.second.bashrc
/root/.inputrc
/root/.less
/root/.vimrc

/etc/default/

/etc/security/limits.conf
/etc/security/limits.conf.ORIG
/etc/security/

/etc/locale.gen
/etc/locale.conf
/etc/vconsole.conf
/etc/environment
/etc/localtime
/etc/fstab
/etc/fstab.ORIG

/etc/host.conf
/etc/hosts
/etc/hosts.ORIG
/etc/hostname

/etc/grub.d/40_custom
/etc/grub.d/
/etc/ntp.conf
/etc/ntp.conf.ORIG
/etc/updatedb.conf
/etc/vimrc
/usr/share/vim/vimfiles/archlinux.vim
/etc/xdg/
/etc/conf.d/
/etc/ssh/







/home/z/build/1nonpkgs/hosts/
/home/z/build/1nonpkgs/rust/
/home/z/build/1nonpkgs/rust.vim/
/home/z/build/1nonpkgs/cargo/

/home/z/build/1nonpkgs/uBlock/
/home/z/build/1nonpkgs/uMatrix/
/home/z/build/1nonpkgs/https-everywhere/
/home/z/build/1nonpkgs/firestorageplus/






/root/.pki
/home/z/.pki
/home/z/.thumbnails
/root/.thumbnails
/root/.cache
/root/.links
/home/z/.cache
/home/z/.cargo
/home/z/.links



/home/z/.../regenkey

/var/log/Xorg.0.log
/home/z/.local/share/xorg/Xorg.0.log
/etc/X11/xorg.conf.d
/etc/X11/

/etc/NetworkManager/dispatcher.d/

/etc/modprobe.d/
/etc/modules-load.d/
/etc/sudoers.d/
/etc/sysctl.d/
/etc/tmpfiles.d/
/etc/udev/rules.d/
/etc/makepkg.conf
/etc/makepkg.conf.ORIG
/etc/pacman.conf.ORIG
/etc/pacman.conf
/etc/resolvconf.conf.ORIG
/etc/resolvconf.conf




/home/z/build/1packages/
/home/z/build/1packages/kernel/linuxgit/makepkg_pacman/linux-git/src/linux-git/.config
/home/z/build/1packages/kernel/linuxgit/makepkg_pacman/linux-git/src/linux-git/.config.old
/home/z/build/1packages/kernel/linuxgit/makepkg_pacman/linux-git/src/linux-git/.version
/home/z/build/1packages/kernel/linuxgit/makepkg_pacman/linux-git/src/linux-git/include/generated/

/home/z/rsync/

/home/z/.gitconfig
/home/z/.xprofile
/home/z/.Xclients
/home/z/.Xmodmap
/home/z/.gitignore
/home/z/.bash_logout
/home/z/.dir_colors
/home/z/.xinitrc
/home/z/.xinitrc.ORIG
/home/z/.gtk-bookmarks
/home/z/.toprc
/home/z/.Xresources
/home/z/.xscreensaver
/home/z/eqmpg123.eq
/etc/pulse/



/etc/iptables/
/etc/systemd/

/home/z/temporary_sandbox_play/

/home/z/.config/psd/
/home/z/.config/libreoffice

/home/z/.gimp-2.8
/home/z/.local/share/mc
/root/.local/share/mc

/etc/cron.weekly/
/etc/pacman.d/gnupg/trustdb.gpg
/etc/pacman.d/gnupg/pubring.gpg
/etc/pacman.d/mirrorlist

/etc/inputrc
/etc/inputrc.ORIG

/etc/dhcpcd.conf
/etc/mkinitcpio.conf
/etc/pamac.conf
/etc/pacman-mirrors.conf

/etc/wgetrc
/etc/lsb-release
/etc/os-release
/etc/manjaro-release
/etc/cpufreq-bench.conf
/etc/gai.conf
/etc/bash.bashrc
/etc/bash.bash_logout
/etc/dnsmasq.conf
/etc/profile
/etc/sensors.d/

/home/z/cscope
/home/z/drate
/home/z/nftdiff
/home/z/backup.nft
/home/z/routerdeny.nft
/home/z/terminfo_linux.src
/home/z/terminfo_xterm.src




/home/z/.config/systemd/user/default.target.wants/
/home/z/.config/systemd/
/home/z/.config/autostart/
/home/z/.config/volumeicon/volumeicon
/usr/lib/systemd/system/cronie.service

/home/z/.local/share/gvfs-metadata
/home/z/build/1packages/chromium-dev-git/fetch_gclient_base/checkout_root/src/.gitattributes
/usbstick/.keep


/home/z/.distcc

/etc/distcc/

/usr/lib/systemd/system-shutdown/debug.sh

/etc/udev/udev.conf

/etc/systemd/system/getty@.service.d/nodisallocate.conf
/usr/lib/systemd/system/cpumax.service
/usr/lib/systemd/system/underclocking_for_nongit_kernels.service
/home/z/Documents/thumbnail

/home/z/manjarred

/etc/ld.so.conf.d/
/etc/ld.so.conf

