# system/Z575/OSes/manjaro2/on_baremetal/filesystem_now/manjaro/usr/lib/systemd/system/cronie.service
[Unit]
Description=Periodic Command Scheduler

[Service]
ExecStart=/usr/bin/crond -n -s
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
